# srcs/ftps/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:vsftpd]
# command=/etc/init.d/vsftpd start
command=/usr/sbin/vsftpd  /etc/vsftpd/vsftpd.conf
autorestart=true
autorestart=true