# srcs/MySQL/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:mariadb]
# command=/etc/init.d/vsftpd start
command=service mariadb start
autorestart=true
autorestart=true