# srcs/MySQL/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:mariadb]
# command=/etc/init.d/vsftpd start
command=/usr/bin/mysqld_safe --datadir='/var/lib/mysql'
autorestart=true
autorestart=true