# srcs/phpMyAdmin/srcs/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true

[program:php-fpm7]
command=/etc/init.d/php-fpm7 start
autostart=true
autorestart=true